# prompter.service
[Service]
ExecStart=/usr/bin/prompter
Restart=always
